# root/etc/services.d/sstvproxy/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc /usr/bin/python3 /config/sstvproxy/sstvProxy.py
